# manishbr/global/global/files/etc/cron.d/winbind-bounce
MAILTO=""
# Restart winbind on the tenth minute, every 6 hours
10 */6 * * * root test -x /usr/bin/systemctl && /usr/bin/systemctl stop winbind || /sbin/service winbind stop; /bin/sleep 5s; /usr/bin/pkill -KILL winbindd; /bin/sleep 2s; test -x /usr/bin/systemctl && /usr/bin/systemctl start winbind || /sbin/service winbind start
